# rtc/supervisor/conf/firestarter-8349.conf
[program:firestarter-8349]
command=/usr/local/bin/luajit /home/overlord/firestarter/server.lua --config=/etc/opt/firestarter/firestarter.conf
stderr_logfile = /var/log/supervisor/firestarter-stderr.log
stdout_logfile = /var/log/supervisor/firestarter-stdout.log
startretries = 2147483647 ; maxint
startsecs = 10