# rtc/supervisor/conf/starfruit-3888.conf
[program:starfruit-3888]
command=/home/asterisk/venv/starfruit/bin/python /home/asterisk/starfruit/server.py --config=/home/asterisk/starfruit/starfruit.conf --port=3888 --mongo_host={{ pillar["starfruit"]["mongo_host"] }} --cache_enabled={{ pillar["starfruit"]["cache_enabled"] }} --memcached_host={{ pillar["starfruit"]["memcached_host"]}}
stderr_logfile = /var/log/supervisor/starfruit-stderr.log
stdout_logfile = /var/log/supervisor/starfruit-stdout.log
startretries = 2147483647 ; maxint
startsecs = 10