# web/supervisor/conf/cebus-8538.conf
[program:cebus-8538]
command=/home/overlord/venv/cebus/bin/python /home/overlord/cebus/server.py --config=/home/overlord/cebus/cebus.conf --port=8538 --mongo_host={{ pillar["cebus"]["mongo_host"] }} --cache_enabled={{ pillar["cebus"]["cache_enabled"] }} --memcached_host={{ pillar["cebus"]["memcached_host"]}} --riak_host={{ pillar["cebus"]["riak_host"] }}
stderr_logfile = /var/log/supervisor/cebus-stderr.log
stdout_logfile = /var/log/supervisor/cebus-stdout.log
startretries = 2147483647 ; maxint
startsecs = 10