# web/supervisor/conf/fiorella-8004.conf
[program:fiorella-8004]
command=/home/overlord/venv/fiorella/bin/python /home/overlord/fiorella/server.py --config=/home/overlord/fiorella/fiorella.conf --port=8004 --mongo_host={{ pillar["fiorella"]["mongo_host"] }} --cache_enabled={{ pillar["fiorella"]["cache_enabled"] }} --memcached_host={{ pillar["fiorella"]["memcached_host"]}}
stderr_logfile = /var/log/supervisor/fiorella-stderr.log
stdout_logfile = /var/log/supervisor/fiorella-stdout.log
startretries = 2147483647 ; maxint
startsecs = 10