# web/supervisor/conf/howler-8001.conf
[program:howler-8001]
command=/home/overlord/venv/howler/bin/python /home/overlord/howler/server.py --config=/home/overlord/howler/howler.conf --port=8001 --mongo_host={{ pillar["howler"]["mongo_host"] }} --cache_enabled={{ pillar["howler"]["cache_enabled"] }} --memcached_host={{ pillar["howler"]["memcached_host"]}} --riak_host={{ pillar["howler"]["riak_host"] }}
stderr_logfile = /var/log/supervisor/howler-stderr.log
stdout_logfile = /var/log/supervisor/howler-stdout.log
startretries = 2147483647 ; maxint
startsecs = 10