# web/supervisor/conf/mango-8002.conf
[program:mango-8002]
command=/home/overlord/venv/mango/bin/python /home/overlord/mango/server.py --config=/home/overlord/mango/mango.conf --port=8002 --mongo_host={{ pillar["mango"]["mongo_host"] }} --sql_host={{ pillar["mango"]["sql_host"] }} --cache_enabled={{ pillar["mango"]["cache_enabled"] }} --memcached_host={{ pillar["mango"]["memcached_host"]}}
stderr_logfile = /var/log/supervisor/mango-stderr.log
stdout_logfile = /var/log/supervisor/mango-stdout.log
startretries = 2147483647 ; maxint
startsecs = 10