# web/supervisor/conf/overlord-8888.conf
[program:overlord-8888]
command=/home/overlord/venv/overlord/bin/python /home/overlord/overlord/server.py --config=/home/overlord/overlord/overlord.conf --port=8888 --mongo_host={{ pillar["overlord"]["mongo_host"] }} --cache_enabled={{ pillar["overlord"]["cache_enabled"] }} --memcached_host={{ pillar["overlord"]["memcached_host"]}} --riak_host={{ pillar["overlord"]["riak_host"] }}
stderr_logfile = /var/log/supervisor/overlord-stderr.log
stdout_logfile = /var/log/supervisor/overlord-stdout.log
startretries = 2147483647 ; maxint
startsecs = 10