# web/supervisor/conf/simon-8539.conf
[program:simon-8539]
command=/home/overlord/venv/simon/bin/python /home/overlord/simon/server.py --config=/home/overlord/simon/simon.conf --port=8539 --mongo_host={{ pillar["simon"]["mongo_host"] }} --cache_enabled={{ pillar["simon"]["cache_enabled"] }} --memcached_host={{ pillar["simon"]["memcached_host"]}} --riak_host={{ pillar["simon"]["riak_host"] }}
stderr_logfile = /var/log/supervisor/simon-stderr.log
stdout_logfile = /var/log/supervisor/simon-stdout.log
startretries = 2147483647 ; maxint
startsecs = 10