# web/supervisor/conf/spider-8003.conf
[program:spider-8003]
command=/home/overlord/venv/spider/bin/python /home/overlord/spider/server.py --config=/home/overlord/spider/spider.conf --port=8003 --mongo_host={{ pillar["spider"]["mongo_host"] }} --sql_host={{ pillar["spider"]["sql_host"] }} --cache_enabled={{ pillar["spider"]["cache_enabled"] }} --memcached_host={{ pillar["spider"]["memcached_host"]}}
stderr_logfile = /var/log/supervisor/spider-stderr.log
stdout_logfile = /var/log/supervisor/spider-stdout.log
startretries = 2147483647 ; maxint
startsecs = 10