# web/supervisor/conf/treehouse-8222.conf
[program:treehouse-8222]
command=/home/overlord/venv/treehouse/bin/python /home/overlord/treehouse/server.py --config=/home/overlord/treehouse/treehouse.conf --port=8222 --mongo_host={{ pillar["treehouse"]["mongo_host"] }} --cache_enabled={{ pillar["treehouse"]["cache_enabled"] }} --memcached_host={{ pillar["treehouse"]["memcached_host"]}} --riak_host={{ pillar["treehouse"]["riak_host"] }}
stderr_logfile = /var/log/supervisor/treehouse-stderr.log
stdout_logfile = /var/log/supervisor/treehouse-stdout.log
startretries = 2147483647 ; maxint
startsecs = 10